# btc-dashboard/deploy/btc-dashboard.service.template
[Unit]
Description=BTC 단타 대시보드 Node.js 서비스
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__/btc-dashboard
EnvironmentFile=__APP_DIR__/btc-dashboard/.env
ExecStart=/usr/bin/node __APP_DIR__/btc-dashboard/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
